# system/netdata-init-d.in
#!/bin/bash
# chkconfig: 345 99 01
# description: startup script

# Source functions
. /etc/rc.d/init.d/functions

DAEMON="netdata"
DAEMON_PATH=@sbindir_POST@
PIDFILE=@localstatedir_POST@/$DAEMON.pid
DAEMONOPTS="-pidfile $PIDFILE"
STOP_TIMEOUT="10"

service_start()
{
	echo "Starting $DAEMON..."
	daemon $DAEMON_PATH/$DAEMON $DAEMONOPTS
	RETVAL=$?
	echo
	return $RETVAL
}

service_stop()
{
	printf "%-50s" "Stopping $DAEMON..."
	killproc -p ${PIDFILE} -d ${STOP_TIMEOUT} $DAEMON
	RETVAL=$?
	if [ $RETVAL -eq 0 ]; then
		rm -f ${PIDFILE}
	fi
	echo
	return $RETVAL
}

service_status()
{
	status -p ${PIDFILE} $DAEMON_PATH/$DAEMON
}

case "$1" in
start)
	service_start
;;
status)
	service_status
;;
stop)
	service_stop
;;
restart)
	service_stop
	service_start
;;
*)
	echo "Usage: $0 {status|start|stop|restart}"
	exit 1
esac
