# system/netdata.service.in
# netdata systemd target

[Unit]
Description=netdata - Real-time performance monitoring
#Documentation=man:netdata
#Documentation=file:///usr/share/doc/netdata/html/index.html
Documentation=https://github.com/firehol/netdata/wiki
After=network.target httpd.service squid.service nfs-server.service mysqld.service named.service postfix.service
Wants=network-online.target
ConditionPathExists=@sysconfig_POST@/netdata/netdata.conf

[Service]
Type=simple
Environment="netdata_LOG_LOCATION=/var/log/netdata/log"
ExecStart=@sbindir_POST@/netdata -D
# Reload is not supported, see https://github.com/firehol/netdata/pull/1277
TimeoutStopSec=60
KillMode=mixed
KillSignal=SIGTERM
OOMScoreAdjust=-900

User=netdata
Group=netdata
Restart=on-abnormal
RestartSec=2s
LimitNOFILE=65536

WorkingDirectory=/tmp

# Hardening

NoNewPrivileges=false
PermissionsStartOnly=true
# CAP_SETGID is required for setgroups()
# CAP_NET_RAW is needed by fping
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_SYS_PTRACE CAP_SETGID CAP_NET_RAW
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
#SystemCallFilter=~@cpu-emulation @debug @keyring @module @mount @obsolete @raw-io @resources

ReadOnlyDirectories=/
ReadWriteDirectories=/proc/self
ReadWriteDirectories=@localstatedir_POST@/lib/netdata
ReadWriteDirectories=@localstatedir_POST@/log/netdata
ReadWriteDirectories=@localstatedir_POST@/cache/netdata

# Access to devices and kernel modules and tunables is required
PrivateDevices=no
ProtectKernelModules=no
ProtectKernelTunables=no

StandardOutput=syslog+console
StandardError=syslog+console

[Install]
WantedBy=multi-user.target
