# system/netdata.service.in
# SPDX-License-Identifier: GPL-3.0-or-later
[Unit]
Description=Real time performance monitoring
#Documentation=man:netdata
#Documentation=file:///usr/share/doc/netdata/html/index.html
Documentation=https://github.com/netdata/netdata/wiki
# append here other services you want netdata to wait for them to start
After=network.target httpd.service squid.service nfs-server.service mysqld.service mysql.service named.service postfix.service chronyd.service

[Service]
Type=simple
User=netdata
Group=netdata
Environment="netdata_LOG_LOCATION=/var/log/netdata/log"
RuntimeDirectory=netdata
RuntimeDirectoryMode=0775
PIDFile=@localstatedir_POST@/run/netdata/netdata.pid
ExecStartPre=/bin/mkdir -p @localstatedir_POST@/cache/netdata
ExecStartPre=/bin/chown -R netdata:netdata @localstatedir_POST@/cache/netdata
ExecStartPre=/bin/mkdir -p @localstatedir_POST@/run/netdata
ExecStartPre=/bin/chown -R netdata:netdata @localstatedir_POST@/run/netdata
ExecStart=@sbindir_POST@/netdata -P @localstatedir_POST@/run/netdata/netdata.pid -D -W set global 'process scheduling policy' 'keep' -W set global 'OOM score' 'keep'
#ExecStopPost=/bin/rm @localstatedir_POST@/run/netdata/netdata.pid
PermissionsStartOnly=true

# Reload is not supported, see https://github.com/firehol/netdata/pull/1277
TimeoutStopSec=60

# restart netdata if it crashes
Restart=on-failure
RestartSec=30
KillMode=mixed
KillSignal=SIGTERM
LimitNOFILE=65536

# The minimum netdata Out-Of-Memory (OOM) score.
# netdata (via [global].OOM score in netdata.conf) can only increase the value set here.
# To decrease it, set the minimum here and set the same or a higher value in netdata.conf.
# Valid values: -1000 (never kill netdata) to 1000 (always kill netdata).
OOMScoreAdjust=1000

# Valid policies: other (the system default) | batch | idle | fifo | rr
# To give netdata the max priority, set CPUSchedulingPolicy=rr and CPUSchedulingPriority=99
CPUSchedulingPolicy=idle

NoNewPrivileges=false
PermissionsStartOnly=true
# CAP_SETGID is required for setgroups()
# CAP_NET_RAW is needed by fping
CapabilityBoundingSet=CAP_DAC_READ_SEARCH CAP_SYS_PTRACE CAP_SETGID CAP_NET_RAW
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full

ReadOnlyDirectories=/
ReadWriteDirectories=/proc/self
ReadWriteDirectories=@localstatedir_POST@/lib/netdata
ReadWriteDirectories=@localstatedir_POST@/log/netdata
ReadWriteDirectories=@localstatedir_POST@/cache/netdata

# Access to devices and kernel modules and tunables is required
PrivateDevices=no
ProtectKernelModules=no
ProtectKernelTunables=no

StandardOutput=syslog+console
StandardError=syslog+console

[Install]
WantedBy=multi-user.target
