# system/netdata-updater.service.in
[Unit]
Description=Daily auto-updates for Netdata
RefuseManualStart=no
RefuseManualStop=yes

[Service]
Type=oneshot
ExecStart=@pkglibexecdir_POST@/netdata-updater.sh
