# neda/neda.service
[Unit]
Description=Neda AI CRM Agent
After=network.target

[Service]
Type=simple
User=neda
Group=neda
WorkingDirectory=/opt/neda

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/neda/google-cloud-sdk/bin:/opt/neda/bin:/opt/netdata/usr/bin"
Environment="HOME=/opt/neda"
Environment="NODE_ENV=production"

# Execute
ExecStart=/opt/neda/bin/neda \
  --agent ./neda.ai \
  --slack \
  --api 8800 \
  --mcp http:8801 \
  --mcp sse:8802 \
  --mcp ws:8803 \
  --openai-completions 8804 \
  --anthropic-completions 8805 \
  --verbose

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=neda
LogNamespace=neda

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
